# Counter Strike/autorestart/server.list
95.128.243.86:27015:/etc/init.d/hldsfun18 restart
95.128.243.86:27020:/etc/init.d/hldsfunru restart